# ntfy-claude-runner/skills/ntfy-claude-runner/resources/com.user.ntfy-claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.ntfy-claude</string>

    <key>ProgramArguments</key>
    <array>
        <!-- uv のパスは `which uv` で確認して書き換えること -->
        <string>__UV_PATH__</string>
        <string>run</string>
        <string>__DAEMON_SCRIPT_PATH__</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NTFY_TOPIC</key>
        <string>__NTFY_TOPIC__</string>
        <key>ZELLIJ_SESSION</key>
        <string>main</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <!-- ログイン時に自動起動 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- クラッシュ時に自動再起動 -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- 再起動間隔（秒） -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- ログ出力先 -->
    <key>StandardOutPath</key>
    <string>__HOME__/.local/share/ntfy-claude/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.local/share/ntfy-claude/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
